# salt/core/systemd/files/prometheus-node-exporter.service
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
User={{ entry.user }}
Group={{ entry.user }}
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/home/{{ entry.user }}/{{ entry.basename }}-{{ entry.version }}.{{ grains.kernel|lower }}-{{ grains.osarch }}/{{ entry.basename }} \
    --web.listen-address :7231 \
{%- if pillar.prometheus.node_exporter.get('smartmon') %}
    --collector.textfile.directory /home/{{ entry.user }}/node-exporter-textfile-directory \
{%- endif %}
    --web.config /home/{{ entry.user }}/config.yaml

[Install]
WantedBy=multi-user.target
